# infrastructure/net-monitor.service
[Unit]
Description=Asynchronous Monitoring Agent Network Guard (Python 3)
Documentation=https://github.com/YourProfile/rpi-network-guard
# Wait strictly for network interfaces and ufw initialization
After=network-online.target
Wants=network-online.target

[Service]
# Least Privilege Architecture
User=pi
Group=pi

# Eliminate asynchronous buffering for clean journalctl logs
Environment=PYTHONUNBUFFERED=1

# Strict path to runtime environment (change to your absolute path)
WorkingDirectory=/home/pi/rpi-network-guard
ExecStart=/usr/bin/python3 /home/pi/rpi-network-guard/src/net_guard/main.py

# Sophisticated self-healing logic
Restart=on-failure
RestartSec=10

[Install]
# Register service during boot process
WantedBy=multi-user.target
